# firmware/files-meshcore-he4025/etc/uci-defaults/95-meshcore-he4025
#!/bin/sh

[ -x /etc/init.d/meshcore ] && /etc/init.d/meshcore enable

exit 0

